# deploy/prayerserver.service
[Unit]
Description=Prayer Server Display
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/prayer-server/prayerserver-display
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=PORT=3000
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
